# deploy/bot.service
[Unit]
Description=Bot Trading Binance
After=network.target

[Service]
User=ubuntu
# IMPORTANTE: Reemplazar /home/ubuntu/proyecto con la ruta real en tu VPS
WorkingDirectory=/home/ubuntu/proyecto
ExecStart=/home/ubuntu/proyecto/venv/bin/python bot_binance.py
Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/proyecto/logs/bot.log
StandardError=append:/home/ubuntu/proyecto/logs/error.log

[Install]
WantedBy=multi-user.target
