# deploy/dashboard.service
[Unit]
Description=Dashboard Streamlit Bot Trading
After=network.target

[Service]
User=ubuntu
# IMPORTANTE: Reemplazar /home/ubuntu/proyecto con la ruta real en tu VPS
WorkingDirectory=/home/ubuntu/proyecto
ExecStart=/home/ubuntu/proyecto/venv/bin/streamlit run dashboard.py --server.port 8501 --server.address 0.0.0.0
Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/proyecto/logs/dashboard.log
StandardError=append:/home/ubuntu/proyecto/logs/dashboard_error.log

[Install]
WantedBy=multi-user.target
